# ipfs.service
[Unit]
Description=IPFSDaemon

[Service]
Type=simple
ExecStart=/home/tvrtko/bin/ipfs daemon

[Install]
WantedBy=default.target
