# startx.service
[Unit]
Description=StartX

[Service]
Type=simple
ExecStart=/usr/bin/startx

[Install]
WantedBy=default.target
